# .chezmoitemplates/osquery/packs/intrusion-detection.conf
{
  "queries": {
    "persistence_launchd": {
      "query": "SELECT path, label, COALESCE(NULLIF(program,''), program_arguments) AS program, run_at_load, keep_alive, start_interval, watch_paths, disabled, username, root_directory, working_directory FROM launchd;",
      "interval": 600,
      "platform": "darwin",
      "description": "New or changed launchd job (LaunchAgent/LaunchDaemon), a freshly planted persistence mechanism."
    },
    "persistence_launchd_overrides": {
      "query": "SELECT label, key, value, uid, path FROM launchd_overrides;",
      "interval": 600,
      "platform": "darwin",
      "description": "Changed launchd Disabled-key override, re-enabling a dormant job or silencing a defense without touching its plist."
    },
    "persistence_startup_items_crontab": {
      "query": "SELECT 'startup_item' AS kind, name, COALESCE(NULLIF(args,''), path) AS command, path, username FROM startup_items UNION ALL SELECT 'crontab' AS kind, event AS name, command, path, '' AS username FROM crontab;",
      "interval": 600,
      "platform": "darwin",
      "description": "New legacy StartupItem or crontab entry, low-noise persistence vectors on modern macOS."
    },
    "listening_ports_non_loopback": {
      "query": "SELECT lp.address, lp.port, lp.protocol, lp.family, p.name, p.path, p.pid, p.uid, p.cmdline FROM listening_ports lp JOIN processes p ON lp.pid = p.pid WHERE lp.port != 0 AND lp.address NOT IN ('127.0.0.1', '::1') AND lp.address NOT LIKE '127.%';",
      "interval": 300,
      "platform": "darwin",
      "description": "New externally-reachable network listener with owning-process attribution, classic backdoor/C2 signal. Excludes ONLY real loopback (127.0.0.0/8, ::1) and the port-0 placeholder; IPv6 link-local (fe80::) is reachable by any host on the same link (RFC 4291), not loopback, so it stays included and the family column carries the address scope."
    },
    "kernel_extensions_new": {
      "query": "SELECT name, version, path, size, linked_against FROM kernel_extensions;",
      "interval": 600,
      "platform": "darwin",
      "description": "Newly loaded kernel extension, ring-0 code / rootkit signal. Volatile refs count omitted so only genuine new kexts surface."
    },
    "system_extensions_new": {
      "query": "SELECT identifier, version, state, category, team, path, bundle_path, mdm_managed FROM system_extensions;",
      "interval": 600,
      "platform": "darwin",
      "description": "New or newly-activated system/DriverKit extension (esp. network_extension that can intercept traffic), persistence/interception signal."
    },
    "suid_bin_unexpected": {
      "query": "SELECT path, username, groupname, permissions FROM suid_bin WHERE path NOT LIKE '/usr/bin/%' AND path NOT LIKE '/usr/sbin/%' AND path NOT LIKE '/usr/libexec/%' AND path NOT LIKE '/System/%' AND path NOT LIKE '/sbin/%' AND path NOT LIKE '/bin/%';",
      "interval": 1800,
      "platform": "darwin",
      "description": "New setuid/setgid binary outside the OS-owned directories, privilege-escalation backdoor signal."
    },
    "recent_logins": {
      "query": "SELECT username, tty, pid, type_name, host, time FROM last WHERE type_name = 'user_process' AND time > (strftime('%s','now') - 86400);",
      "interval": 600,
      "platform": "darwin",
      "description": "Interactive login in the last 24h, with remote host where present. Unexpected user/tty/host is an account-compromise signal."
    }
  }
}
