# Library/LaunchAgents/com.webdavis.osquery-posture-watch.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.webdavis.osquery-posture-watch</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/bash</string>
    <string>-l</string>
    <string>{{ .chezmoi.homeDir }}/.local/bin/osquery-posture-watch.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>60</integer>
  <key>StandardOutPath</key>
  <string>{{ .chezmoi.homeDir }}/.local/log/osquery/posture-watch.log</string>
  <key>StandardErrorPath</key>
  <string>{{ .chezmoi.homeDir }}/.local/log/osquery/posture-watch.log</string>
</dict>
</plist>
