# Library/LaunchAgents/com.webdavis.atuin-daemon.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.webdavis.atuin-daemon</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/atuin</string>
    <string>daemon</string>
    <string>start</string>
    <!-- --force: clean up stale socket file from a prior crashed daemon
         before binding. Without this, an unclean exit leaves the Unix
         socket file behind and the next start fails with EADDRINUSE,
         causing launchd to enter a tight restart loop (~6000 attempts
         observed in production). Safe for single-user LaunchAgent: at
         most one daemon instance ever exists. -->
    <string>--force</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>{{ .chezmoi.homeDir }}/.local/log/atuin-daemon.log</string>
  <key>StandardErrorPath</key>
  <string>{{ .chezmoi.homeDir }}/.local/log/atuin-daemon.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>{{ .chezmoi.homeDir }}</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
