# Library/LaunchAgents/com.webdavis.openclaw-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.webdavis.openclaw-watchdog</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/stephen/.local/bin/gateway-watchdog.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/stephen/.openclaw/logs/watchdog.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/stephen/.openclaw/logs/watchdog.err.log</string>
  </dict>
</plist>
