# Library/LaunchAgents/com.webdavis.update-skills.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.webdavis.update-skills</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/bash</string>
    <string>{{ .chezmoi.homeDir }}/.local/bin/update-skills.sh</string>
    <!-- Marks this as a scheduled (not manual) run: only a scheduled run with no
         later slot remaining this week claims retry-budget exhaustion. -->
    <string>--scheduled</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>{{ .chezmoi.homeDir }}/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>{{ .chezmoi.homeDir }}</string>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <!-- Four Monday retry slots (04:00/08:00/12:00/16:00). A slot that finds an
       interactive harness session defers and leaves no weekly success stamp, so
       the next slot retries; once a slot succeeds, the rest no-op on the stamp.
       The 16:00 slot is the LAST (a deferral there alerts loudly). Keep the last
       hour in sync with UPDATE_SKILLS_LAST_SLOT_HOUR in update-skills.sh. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>4</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>8</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>16</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>
  <key>StandardOutPath</key>
  <string>{{ .chezmoi.homeDir }}/.local/log/skills/update-skills.log</string>
  <key>StandardErrorPath</key>
  <string>{{ .chezmoi.homeDir }}/.local/log/skills/update-skills.log</string>
</dict>
</plist>
